# conf_templates/gunicorn.conf
description "%(PROJECT_NAME)s webserver (backend)"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec %(DEPLOY_DIR)s/conf/gunicorn.sh
